# crontab

2 16 * * 1-5 python3 /opt/poll/send_poll.py > /dev/null 2>&1

